# rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv sh
set -eu

exec /command/s6-setuidgid atbox php-fpm --nodaemonize -O --fpm-config /etc/php/8.3/fpm/php-fpm.conf
