# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv sh
set -eu

exec /command/s6-setuidgid atbox nginx -g 'daemon off;'
