# rootfs/etc/services.d/php-fpm/run
#!/command/with-contenv sh
set -eu

exec /command/s6-setuidgid "${ATBOX_RUNTIME_USER:-atbox}" php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.d/atom.conf
