# rootfs/etc/services.d/php-fpm/run
#!/command/with-contenv sh
set -eu

exec php-fpm --nodaemonize --allow-to-run-as-root

